# scripts/supervisor-mentorde.conf
; MentorDE — Supervisor Queue Worker Yapılandırması
;
; Kurulum (Hostinger VPS):
;   sudo cp scripts/supervisor-mentorde.conf /etc/supervisor/conf.d/mentorde.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start mentorde:*
;
; .env ayarları (production):
;   QUEUE_CONNECTION=database   (Redis yoksa)
;   QUEUE_CONNECTION=redis      (Redis varsa — daha hızlı)
;
; Database queue için migration: php artisan queue:table && php artisan migrate

; ────────────────────────────────────────────────
; 1. Ana Queue Worker (default, notifications, emails)
;    QUEUE_CONNECTION=database kullanıyorsa "database" yaz
;    QUEUE_CONNECTION=redis    kullanıyorsa "redis" yaz
; ────────────────────────────────────────────────
[program:mentorde-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mentorde/artisan queue:work database --queue=default,notifications,emails --sleep=3 --tries=3 --max-time=3600 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/mentorde/storage/logs/queue-default.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=60

; ────────────────────────────────────────────────
; 2. Yüksek Öncelik Queue (task_due_reminder, escalation)
; ────────────────────────────────────────────────
[program:mentorde-queue-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mentorde/artisan queue:work database --queue=high,task_due_reminder,escalation --sleep=3 --tries=3 --max-time=3600 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/mentorde/storage/logs/queue-high.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stopwaitsecs=60

; ────────────────────────────────────────────────
; 3. Laravel Scheduler (cron yerine supervisor alternatifi)
;    Not: Cron tercih ediliyorsa bu program'ı disable et.
;    Cron kullanımı: * * * * * www-data cd /var/www/mentorde && php artisan schedule:run >> /dev/null 2>&1
; ────────────────────────────────────────────────
[program:mentorde-scheduler]
process_name=%(program_name)s
command=php /var/www/mentorde/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/mentorde/storage/logs/scheduler.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stopwaitsecs=30

; ────────────────────────────────────────────────
; Grup tanımı
; ────────────────────────────────────────────────
[group:mentorde]
programs=mentorde-queue-default,mentorde-queue-high,mentorde-scheduler
priority=999
